# scripts/com.wkerr.knowledge-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wkerr.knowledge-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/wkerr/development/knowledge-management/scripts/run-with-env.sh</string>
        <string>/Users/wkerr/.local/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>/Users/wkerr/development/knowledge-management</string>
        <string>gunicorn</string>
        <string>--workers</string>
        <string>3</string>
        <string>--bind</string>
        <string>0.0.0.0:5001</string>
        <string>knowledge.server:app</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/wkerr/development/knowledge-management</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/knowledge-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/knowledge-server.err</string>
</dict>
</plist>
